# inwx-update.timer
[Unit]
Description=Run inwx-update every five minutes and on boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
